# dist/general-ludd.service
[Unit]
Description=General Ludd Agent Daemon
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=general-ludd
EnvironmentFile=-/etc/general-ludd/env
ExecStart=/usr/local/bin/gludd daemon --host 127.0.0.1 --port 8000
Restart=on-failure
RestartSec=5s
NoNewPrivileges=true
ProtectSystem=strict
PrivateTmp=true

[Install]
WantedBy=multi-user.target
