# dist/hottentot-agent-0.1.0-Darwin-arm64/hottentot.service
[Unit]
Description=Hottentot Agent Daemon
After=network.target
Wants=network-online.target

[Service]
Type=simple
ExecStart=/usr/local/bin/hottentot daemon --host 0.0.0.0 --port 8000 --log-level info
Restart=on-failure
RestartSec=5
WorkingDirectory=/etc/hottentot
Environment=PYTHONUNBUFFERED=1

# Security hardening
NoNewPrivileges=true
ProtectSystem=strict
ProtectHome=true
ReadWritePaths=/var/log/hottentot /var/lib/hottentot
PrivateTmp=true

[Install]
WantedBy=multi-user.target
